# init/synapse.conf
#!upstart

description "synapse"
author "er0k"

start on startup
stop on shutdown
respawn

script
    echo $$ > /var/run/synapse.pid
    sudo -u /opt/synapse/server.py >> /var/log/synapse
end script