# init/synapse.conf
#!upstart

description "synapse"
author "er0k"

start on startup
stop on shutdown
respawn

exec /opt/synapse/server.py
